# files/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=service php7.0-fpm start
autostart=true
autorestart=true